# deploy/spielesammlung.service
[Unit]
Description=Spielesammlung Flask App
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/spielesammlung
Environment=SECRET_KEY=CHANGE_ME
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/spielesammlung/venv/bin/gunicorn \
    --worker-class eventlet \
    -w 1 \
    --bind 127.0.0.1:5000 \
    --access-logfile /var/log/spielesammlung/access.log \
    --error-logfile /var/log/spielesammlung/error.log \
    "app:create_app()"
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
